# roles/webserver/templates/chronograph.service.j2
[Unit]
Description=Chronograph Webserver

[Service]
ExecStart=java -jar /opt/chronograph.jar -sf /etc/chronograph/config.edn

[Install]
WantedBy=multi-user.target
